# data/laniakea-webswview.service.in
[Unit]
Description=Laniakea Web Software View
After=network.target
ConditionPathExists=/var/lib/laniakea/webswview/config.cfg

[Service]
Restart=on-failure

#ProtectSystem=full
#ProtectHome=true

User=lkweb
Group=www-data
WorkingDirectory=@LKSWWEB_INSTALL_DIR@

RuntimeDirectory=laniakea-webswview
ExecStart=/usr/bin/uwsgi --ini lkwebswview.ini

[Install]
WantedBy=multi-user.target
