# data/laniakea-scheduler.service.in
[Unit]
Description=Laniakea Maintenance Daemon
Requires=postgresql.service
After=network.target
ConditionPathExists=/etc/laniakea/base-config.toml

[Service]
Type=notify
Restart=on-failure
TimeoutStartSec=10

User=lkmaster
Group=lkmaster
ExecStart=@SCHEDULER_INSTALL_BIN@

[Install]
WantedBy=multi-user.target
