# data/laniakea-webdash.service.in
[Unit]
Description=Laniakea Web Dashboard
After=syslog.target network.target
ConditionPathExists=/var/lib/laniakea/webdash/config.cfg

[Service]
Restart=on-failure

#ProtectSystem=full
#ProtectHome=true

User=lkweb
Group=www-data
WorkingDirectory=@LKWEB_INSTALL_DIR@

RuntimeDirectory=lkweb
ExecStart=/usr/bin/uwsgi --ini lkwebdash.ini

[Install]
WantedBy=multi-user.target
