# data/laniakea-mirk.service.in
[Unit]
Description=Laniakea Matrix Bot
After=network.target
ConditionPathExists=/etc/laniakea/mirk.toml

[Service]
Restart=on-failure
RestartSec=15

#ProtectSystem=full
#ProtectHome=true

User=lkweb
Group=nogroup
ExecStart=@MIRK_INSTALL_BIN@

[Install]
WantedBy=multi-user.target
