# data/laniakea-lighthouse.service.in
[Unit]
Description=Laniakea Job Server
Requires=postgresql.service
After=syslog.target network.target
ConditionPathExists=/etc/laniakea/base-config.toml

[Service]
Type=notify
Restart=on-failure
ProtectSystem=full
User=lklighthouse
Group=nogroup
ExecStart=@LIGHTHOUSE_INSTALL_BIN@

[Install]
WantedBy=multi-user.target
